# deploy/systemd/daxson-serve.service
[Unit]
Description=Daxson Tunnel Server (device-auth + dashboard mode)
Documentation=https://github.com/daxson/tunnel
After=network-online.target nginx.service
Wants=network-online.target
# Ensure nginx is running before we start (probe_upstream depends on it)
BindsTo=nginx.service

[Service]
Type=simple
User=daxson
Group=daxson

# daxson serve reads ~/.daxson/server-identity.json and ~/.daxson/registry.json
# (relative to the running user's home) unless explicit paths are set in server.yaml.
ExecStart=/usr/local/bin/daxson serve --config /etc/daxson/server.yaml
ExecReload=/bin/kill -HUP $MAINPID

# Restart aggressively — availability is critical.
Restart=on-failure
RestartSec=5s
StartLimitBurst=10
StartLimitIntervalSec=60

# File descriptor and process limits.
LimitNOFILE=131072
LimitNPROC=8192

# Security hardening.
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=false
# Allow writes to the daxson user home (for identity + registry auto-generation).
ReadWritePaths=/home/daxson
# Allow reads from /etc/daxson (config + certs).
ReadOnlyPaths=/etc/daxson /etc/letsencrypt
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

StandardOutput=journal
StandardError=journal
SyslogIdentifier=daxson-serve

[Install]
WantedBy=multi-user.target
